# ansible/roles/common/files/update-motd.cron
# Delay a random time between 0 and 360 minutes (6 hours) to avoid having all
# instances update at the same time.
RANDOM_DELAY=360
@daily root /sbin/start --quiet update-motd